# files/etc/uci-defaults/99-m2-aurora-default-home
#!/bin/sh

echo "==== Set M2 Aurora as LuCI default home ===="

uci -q set luci.main.mediaurlbase='/luci-static/aurora'
uci -q commit luci

INDEX_LUA="/usr/lib/lua/luci/controller/admin/index.lua"

if [ -f "$INDEX_LUA" ]; then
    sed -i 's/alias("admin", "status", "overview")/alias("admin", "status", "m2aurora")/g' "$INDEX_LUA"
    sed -i "s/alias('admin', 'status', 'overview')/alias('admin', 'status', 'm2aurora')/g" "$INDEX_LUA"
    sed -i 's/firstchild()/alias("admin", "status", "m2aurora")/g' "$INDEX_LUA"
fi

MENU_DIR="/usr/share/luci/menu.d"
BASE_JSON="$MENU_DIR/luci-base.json"

mkdir -p "$MENU_DIR"

cat > "$MENU_DIR/m2-aurora-home.json" <<'JSON'
{
  "admin/status/m2aurora": {
    "title": "M2 Aurora 首页",
    "order": -100,
    "action": {
      "type": "template",
      "path": "m2more/aurora"
    },
    "depends": {
      "fs": {
        "/usr/lib/lua/luci/view/m2more/aurora.htm": "file"
      }
    }
  }
}
JSON

if [ -f "$BASE_JSON" ]; then
    sed -i 's#admin/status/overview#admin/status/m2aurora#g' "$BASE_JSON"
fi

rm -rf /tmp/luci-indexcache /tmp/luci-modulecache /tmp/luci-indexcache.* /tmp/luci-modulecache.* /tmp/luci-*cache* 2>/dev/null || true

exit 0
